# install/opt/local/etc/profile.d/10-paths.sh
#!/bin/sh

export PATH="/opt/local/bin:/opt/local/sbin:$PATH" # make sure user would access macports first instead of system
export PATH="$HOME/.local/bin:$PATH"
export PATH="$HOME/.local/share/krew/bin:$PATH"
