# Library/LaunchAgents/50.environment.program.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
		<key>Label</key>
		<string>environment.program</string>
		<key>ProgramArguments</key>
		<array>
				<string>/bin/sh</string>
				<string>-c</string>
				<string>
					/bin/launchctl setenv DOCKER_BUILDKIT "1"
					/bin/launchctl setenv GOPROXY "https://goproxy.io,direct"
					/bin/launchctl setenv SBT_OPTS "-Dsbt.override.build.repos=true"
				</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
		<key>StandardErrorPath</key>
		<string>/dev/null</string>
		<key>StandardOutPath</key>
		<string>/dev/null</string>
</dict>
</plist>
