# Library/LaunchAgents/70.environment.program.xdg.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
		<key>Label</key>
		<string>environment.program.xdg</string>
		<key>ProgramArguments</key>
		<array>
				<string>/bin/sh</string>
				<string>-c</string>
				<string>
					/bin/launchctl setenv DOCKER_CONFIG "$HOME/.config/docker"
					/bin/launchctl setenv MACHINE_STORAGE_PATH "$HOME/.local/share/docker-machine"
					/bin/launchctl setenv GRADLE_USER_HOME "$HOME/.local/share/gradle"
					/bin/launchctl setenv NPM_CONFIG_USERCONFIG "$HOME/.config/npm/npmrc"
					/bin/launchctl setenv NODE_REPL_HISTORY "$HOME/.local/share/node_repl_history"
					/bin/launchctl setenv WORKON_HOME "$HOME/.local/share/virtualenvs"
					/bin/launchctl setenv HISTFILE "$HOME/.local/state/bash/history"
					/bin/launchctl setenv GNUPGHOME "$HOME/.local/share/gnupg"
					/bin/launchctl setenv INPUTRC "$HOME/.config/readline/inputrc"
					/bin/launchctl setenv RIPGREP_CONFIG_PATH "$HOME/.config/ripgrep/config"
					/bin/launchctl setenv GDBHISTFILE "$HOME/.local/share/gdb/history"
					/bin/launchctl setenv RUSTUP_HOME "$HOME/.local/share/rustup"
					/bin/launchctl setenv CARGO_HOME "$HOME/.local/share/cargo"
					/bin/launchctl setenv GOPATH "$HOME/.local/share/go"
					/bin/launchctl setenv PYLINTHOME "$HOME/.cache/pylint"
					/bin/launchctl setenv PASSWORD_STORE_DIR "$HOME/.local/share/pass"
					/bin/launchctl setenv K9SCONFIG "$HOME/.config/k9s"
					/bin/launchctl setenv VAGRANT_HOME "$HOME/.local/share/vagrant"
					/bin/launchctl setenv VAGRANT_ALIAS_FILE "$HOME/.local/share/vagrant/aliases"
					/bin/launchctl setenv TEXMFHOME "$HOME/.local/share/texmf"
					/bin/launchctl setenv TEXMFVAR "$HOME/.cache/texlive/texmf-var"
					/bin/launchctl setenv TEXMFCONFIG "$HOME/.config/texlive/texmf-config"
				</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
		<key>StandardErrorPath</key>
		<string>/dev/null</string>
		<key>StandardOutPath</key>
		<string>/dev/null</string>
</dict>
</plist>
