# Library/LaunchAgents/90.environment.pass.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
		<key>Label</key>
		<string>environment.pass</string>
		<key>ProgramArguments</key>
		<array>
				<string>/bin/sh</string>
				<string>-c</string>
				<string>
					/bin/launchctl setenv PASSWORD_STORE_CHARACTER_SET "a-zA-Z0-9~!@#$%^&*()-_=+[]{};:,.<>?"
					/bin/launchctl setenv PASSWORD_STORE_GENERATED_LENGTH "40"
				</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
		<key>StandardErrorPath</key>
		<string>/dev/null</string>
		<key>StandardOutPath</key>
		<string>/dev/null</string>
</dict>
</plist>
