# .local/share/brew/etc/profile.d/10-paths.sh
#!/bin/sh

export PATH="$HOME/.local/share/krew/bin:$PATH"
export PATH="$HOME/.local/share/brew/bin:$HOME/.local/share/brew/sbin:$PATH"
export PATH="$HOME/.local/bin:$PATH"
