# .local/etc/profile.d/00-system.sh
export EDITOR="nvim"
export VISUAL="bat"
export PAGER="less"
export DIFFPROG="delta"
