# .local/etc/profile.d/10-path.sh
#!/bin/sh

export PATH="$HOME/.krew/bin:$PATH"
export PATH="$HOME/.npm/bin:$PATH"
export PATH="$HOME/Library/Python/3.10/bin:$PATH"
export PATH="/opt/local/bin:/opt/local/sbin:$PATH"
export PATH="$HOME/.local/bin:$PATH"
