# .local/etc/profile.d/30-toolsets.sh
#!/bin/sh

[[ -f ~/.local/share/cargo/env ]] && source ~/.local/share/cargo/env
