# .local/etc/profile.d/30-toolset.sh
#!/bin/sh

[[ -f ~/.cargo/env ]] && source ~/.cargo/env
